# debian/postinst
#!/bin/sh

#rc-update del powervr default
rc-update add powervr sysinit

ln -s /lib/ld-linux-armhf.so.3 /lib/ld-linux.so.3
