# debian/postinst
#!/bin/sh

#rc-update del powervr default
rc-update add powervr sysinit

if [ "$DEB_HOST_ARCH" = armhf ]; then
	ln -svnf /lib/ld-linux-armhf.so.3 /lib/ld-linux.so.3
fi
