# debian/postinst
#!/bin/sh

#rc-update del powervr default
rc-update add powervr sysinit
